# templates/vector.service.j2
[Unit]
Description=Vector
Documentation=https://vector.dev
After=network-online.target
Requires=network-online.target

[Service]
User=root
Group=root
ExecStart={{ vector_install_dir }}/bin/vector --config {{ vector_config_dir }}/vector.yml
Restart=always

[Install]
WantedBy=multi-user.target
